# docker/supervisor/supervisor.conf
[supervisord]
nodaemon=true

[program:laravel-worker]
command=php /var/www/artisan queue:work
process_name=%(program_name)s_%(process_num)02d
numprocs=1
priority=999
autostart=true
autorestart=true
startsecs=1
startretries=3
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/worker.log
